# web_project/background/supervisord.conf
[unix_http_server]
file = /var/run/back_supervisord.sock

[supervisorctl]
serverurl = unix:///var/run/back_supervisord.sock

[supervisord]
logfile=/opt/logs/supervisord/back_supervisord.log ; (main log file;default $CWD/supervisord.log)
logfile_maxbytes=50MB       ; (max main logfile bytes b4 rotation;default 50MB)
logfile_backups=10          ; (num of main logfile rotation backups;default 10)
loglevel=warn               ; (log level;default info; others: debug,warn,trace)
pidfile=/var/run/back_supervisord.pid ; (supervisord pidfile;default supervisord.pid)
minfds=1024                 ; (min. avail startup file descriptors;default 1024)
minprocs=200                ; (min. avail process descriptors;default 200)
user=root                 ; (default is current user, required if root)
childlogdir=/opt/logs/supervisord/            ; ('AUTO' child log dir, default $TEMP)


[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface


[program:clean_cache]
command=/opt/python2.7/bin/python /opt/app/python/m-cms-new/background/clean_cache/do_clean_cache.py
process_name=%(program_name)s_%(process_num)02d
numprocs=8
autostart=true


[program:clean_url_cache]
command=/opt/python2.7/bin/python /opt/app/python/m-cms-new/background/clean_cache/do_url_clean.py
process_name=%(program_name)s_%(process_num)02d
numprocs=12
autostart=true


[program:clean_box]
command=/opt/python2.7/bin/python /opt/app/python/m-cms-new/background/clean_box/do_clean_box.py
process_name=%(program_name)s_%(process_num)02d
numprocs=8
autostart=true


[program:async_channel_job]
command=/opt/python2.7/bin/python /opt/app/python/m-cms-new/background/sync_data/async_channel_job.py
process_name=%(program_name)s_%(process_num)02d
numprocs=100
autostart=true

